# 1.0/src/ssh.sv.conf
[program:ssh]
command=/usr/local/bin/start_ssh.sh
autostart=true
autorestart=true
